# restful-exercises/goproject.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:myserver]
command = /home/champuser/CapstoneDemo/restful-exercises/bin/romanserver
autostart=true
autorestart=true
redirect_stderr=true
